# vagrantland/whats-a-consul-cluster/files/nginx/etc/network/if-pre-up.d/iptables
#!/bin/bash
#
# Loop through all files in /etc/iptables.d directory

for f in /etc/iptables.d/*.rules
do
  /sbin/iptables-restore < $f
done
